# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/share/nginx/sbin/nginx -g "daemon off;"

[program:php5-fpm]
command=/usr/sbin/php5-fpm
